# stowage/work-machine/Library/LaunchAgents/clean-kube-config-files.plist
  <?xml version="1.0" encoding="UTF-8"?>
  <!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
  <plist version="1.0">
  <dict>
      <key>Label</key>
      <string>clean-kube-config-files</string>

      <key>ProgramArguments</key>
      <array>
          <!-- somewhat fragile -->
          <string>/opt/homebrew/bin/bash</string>
          <string>-c</string>
          <string>source $HOME/.bash_profile &amp;&amp; $HOME/config/non-stowage/clean-kube-config-files.sh</string>
      </array>

      <key>StartCalendarInterval</key>
      <dict>
          <key>Hour</key>
          <integer>18</integer>
          <key>Minute</key>
          <integer>0</integer>
      </dict>

      <key>StandardOutPath</key>
      <string>/tmp/clean-kube-config-files.stdout.log</string>

      <key>StandardErrorPath</key>
      <string>/tmp/clean-kube-config-files.stderr.log</string>

      <key>RunAtLoad</key>
      <false/>
  </dict>
  </plist>
